# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/google/bin/uwsgi --ini /srv/google/.config/uwsgi/uwsgi-app.ini
